# etc/ddns-cpp.service
[Unit]
Description=ddns-cpp
After=network.target
[Service]
ExecStart=/usr/local/bin/ddns-cpp
Restart=on-failure
[Install]
WantedBy=multi-user.target
